# Task6/files/cpu.service
[Unit]  
Description=CPU load
After=network.target  
StartLimitIntervalSec=0

[Service]  
Type=simple  
Restart=always  
RestartSec=1  
User=ubuntu  
ExecStart=/home/ubuntu/loadCPU.sh
  
[Install]  
WantedBy=multi-user.target